# overlay/etc/local.d/usr.start
[ ! -d "/overlay/usr" ] && mkdir /overlay/{usr,usrw}
mount -t overlay overlay -o lowerdir=/usr,upperdir=/overlay/usr,workdir=/overlay/usrw,ro /usr